# roles/initrd/files/etc/network/if-up.d/showme
#!/bin/sh

echo "if-up trying: $*"
env > /if-up.env