# roles/initrd/files/etc/network/if-pre-up.d/showme
#!/bin/sh

echo "if-pre-up trying: $*"
env > /if-pre-up.env